# server/quantum-flow.service
[Unit]
Description=QUANTUM FLOW v2.1 — KOSPI/KOSDAQ 추세추종 자동매매 시스템
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/quantum-flow
EnvironmentFile=/home/ubuntu/quantum-flow/.env
ExecStart=/home/ubuntu/quantum-flow/.venv/bin/python main.py

# ── 자동 재시작 설정 ──────────────────────────────────────────────
# 크래시 후 30초 대기 후 재시작, 5분 내 5회 실패 시 포기
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# ── 로그 설정 (journalctl -u quantum-flow 로 확인) ────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-flow

# ── 메모리 보호 (4GB RAM 기준, 3GB 상한) ─────────────────────────
MemoryMax=3G
MemoryHigh=2.5G

# ── 프로세스 보호 (Lightsail 환경) ───────────────────────────────
# 파일 디스크립터 제한 완화 (웹소켓 연결 다수)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
